# assets/config/supervisor/5-remove-singleton.conf
[program:Chromium-Oneshot-Prepare]
command=bash -c 'rm -f /root/.config/chromium/SingletonLock; rm -f /root/.config/chromium/SingletonLock;'
priority=5
stderr_logfile=/var/log/chromium.err.log
